# ebusd/rootfs/etc/services.d/ebusd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: eBUSd
# Runs ebusd
# ==============================================================================


bashio::log.info "Starting ebusd..."

if [ -d /data/ebusd-configuration ]; 
then 
    if [ -d /etc/ebusd ]; then sudo mv /etc/ebusd /etc/ebusd.old; fi
    sudo ln -s /data/ebusd-configuration/ebusd-2.1.x/de /etc/ebusd
fi

exec /usr/bin/ebusd \
    --scanconfig \
    --device=$(bashio::config 'device') \
    --port=$(bashio::addon.port '8888/tcp') \
    --foreground  \
    --mqtthost=$(bashio::services mqtt "host") \
    --mqttport=$(bashio::services mqtt "port") \
    --mqttuser=$(bashio::services mqtt "username") \
    --mqttpass=$(bashio::services mqtt "password")