# ebusd/rootfs/etc/services.d/ebusd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: ebusd
# Runs ebusd script
# ==============================================================================


bashio::log.info "Starting ebusd..."

exec /usr/bin/ebusd --scanconfig --device=$(bashio::config 'device') --port=$(bashio::addon.port '8888/tcp') --foreground
