# ebusd/rootfs/etc/services.d/ebusd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: eBUSd
# Runs ebusd
# ==============================================================================




EBUSDCONFIGPATH="http://ebusd.eu/config/";
if [ -d /config/ebusd-configuration ]; 
then 
    EBUSDCONFIGPATH="/config/ebusd-configuration";
fi

bashio::log.info "Using eBUSd configuration path: '${EBUSDCONFIGPATH}'"
bashio::log.info "Starting ebusd..."
exec /usr/bin/ebusd \
    --foreground  \
    --scanconfig \
    --configpath "${EBUSDCONFIGPATH}" \
    --device=$(bashio::config 'device') \
    --port=$(bashio::addon.port '8888/tcp') \
    --mqtthost=$(bashio::services mqtt "host") \
    --mqttport=$(bashio::services mqtt "port") \
    --mqttuser=$(bashio::services mqtt "username") \
    --mqttpass=$(bashio::services mqtt "password")